# etc/supervisord.conf
[supervisord]
nodaemon=true

[program:vnc]
command=sh /home/debian/vncserver.sh
autostart=true
user=debian
redirect_stderr=true

